# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#echo netdev > /sys/class/leds/longsys\:blue\:led/trigger
#echo wl0.1 > /sys/class/leds/longsys\:blue\:led/device_name
#echo link tx rx > /sys/class/leds/longsys\:blue\:led/mode

init_network.sh

updatefw &

set_to_lan &

/etc/init.d/samba restart
exit 0
